# automation/zion-error-monitor.service
[Unit]
Description=Zion Console Error Monitor Daemon
Documentation=https://github.com/Zion-Holdings/zion.app
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/zion.app
ExecStart=/usr/bin/node /var/www/zion.app/automation/continuous-error-monitor.js --daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zion-error-monitor

# Environment variables
Environment=NODE_ENV=production
Environment=ZION_INTEGRATION=true

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/zion.app/reports

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Process management
PIDFile=/var/www/zion.app/reports/monitor.pid

[Install]
WantedBy=multi-user.target
